# scripts/discord-bot.service
[Unit]
Description=Discord Zomboid Bot
After=network.target

[Service]
Type=simple
User=pzuser
WorkingDirectory=/opt/discord-bot
Environment="PATH=/opt/discord-bot/venv/bin"
ExecStart=/opt/discord-bot/venv/bin/python main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target